# systemd/caracas.target
# vi: se ft=systemd:

[Unit]
Description=All caracas services
Requires=multi-user.target
After=multi-user.target
