# systemd/caracas.target
# vi: se ft=systemd:

[Unit]
Description=All caracas services
Requires=multi-user.target
Conflicts=rescue.target
AllowIsolate=yes
